# messenger.service
[Unit]
Description=Messenger — agent-to-agent messaging daemon
After=network.target

[Service]
Type=simple
User=richard
WorkingDirectory=/home/richard/projects/messenger
ExecStart=/usr/bin/java -jar /home/richard/projects/messenger/target/messenger.jar /home/richard/projects/messenger/config.json
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Pass OpenBrain key from environment
EnvironmentFile=-/home/richard/projects/messenger/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/richard

# v1.2 progress beats (issue #17): per-thread log files live under
# /run/messenger/progress (== /var/run/messenger/progress). systemd creates
# this directory with mode 0755 owned by the service user on each start, and
# wipes it on stop, so the daemon never has to manage permissions itself.
RuntimeDirectory=messenger/progress
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
